# openwrt/files/etc/crontabs/root
*/5 * * * *	sh /etc/SecureAP/periodic.sh


